# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:frontpage]
command=uv --project /app/apps/frontpage run gunicorn -k uvicorn.workers.UvicornH11Worker -b 127.0.0.1:5000 --chdir /app/apps frontpage:app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:hello]
command=gunicorn -k eventlet -b 127.0.0.1:5001 --chdir /app/apps hello:app
directory=/app/apps/hello
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:parking]
command=gunicorn -k eventlet -b 127.0.0.1:5002 --chdir /app/apps parking:app
directory=/app/apps/parking
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:nationofpositivity]
command=/strip-env.sh NATIONOFPOSITIVITY_ -- gunicorn -k eventlet -b 127.0.0.1:5004 --chdir /app/apps nationofpositivity.shop:server
directory=/app/apps/nationofpositivity
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:homemadebycvg_com]
command=gunicorn -k eventlet -b 127.0.0.1:5005 --chdir /app/apps homemadebycvg:app
directory=/app/apps
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:getijden]
command=gunicorn -k eventlet -b 127.0.0.1:5006 --chdir /app/apps/getijden web:app
directory=/app/apps/getijden
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:letmelearn]
command=/strip-env.sh LETMELEARN_ -- gunicorn -k eventlet -b 127.0.0.1:5007 --chdir /app/apps/letmelearn letmelearn.web:server
directory=/app/apps/letmelearn
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:baseweb_demo]
command=/strip-env.sh BASEWEB_DEMO_ -- gunicorn -k eventlet -b 127.0.0.1:5008 --chdir /app/apps baseweb-demo:server
directory=/app/apps
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:howifeel]
command=/strip-env.sh HOWIFEEL_ -- gunicorn -k eventlet -b 127.0.0.1:5009 --chdir /app/apps/howifeel howifeel:app
directory=/app/apps/howifeel
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10

[program:oatk]
command=/strip-env.sh OATK_DEMO_ -- gunicorn -k eventlet -b 127.0.0.1:5010 --chdir /app/apps/oatk examples.google.app:server
directory=/app/apps/oatk
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
priority=10
